# king/etc/supervisor/conf.d/myconf.ini
[program:myweb_dj]
command=pipenv run gunicorn myweb_dj.wsgi -w 2 -k gthread -b 127.0.0.1:8000
directory=/home/king/myweb_dj
autostart=true
autorestart=unexpected
user=root
stdout_logfile=/home/king/etc/supervisor/var/log/hellodjango-blog-tutorial-stdout.log
stderr_logfile=/home/king/etc/supervisor/var/log/hellodjango-blog-tutorial-stderr.log
